# supervisor.conf
[program:must-app]
command = python3 /opt/matter/src/serve.py

[program:nginx]
command = /usr/sbin/nginx